# .config/systemd/user/battery-checker.service
[Unit]
Description=Checks battery level and reports if its too low

[Service]
ExecStart=/home/kobu/.local/bin/battery_checker
